# OpenWRT/files/etc/crontabs/root
0 0 * * * /usr/sbin/ntpdate 0.openwrt.pool.ntp.org > /tmp/runlog.log 2>&1
